# deb-systemd/DEBIAN/postinst
#!/bin/sh -e

deb-systemd-helper enable airtop-fpsvc.service || :
deb-systemd-invoke start airtop-fpsvc.service || :
